# machines/deb-lap/scripts/clipboard-listener.service
[Unit]
Description=Clipboard listener for SSH reverse tunnel
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/clipboard-listener.sh
Restart=always

[Install]
WantedBy=graphical-session.target
