# machines/deb-lap/scripts/clipboard-listener.service
[Unit]
Description=Clipboard listener for SSH reverse tunnel

[Service]
ExecStart=%h/.local/bin/clipboard-listener.sh
Restart=always

[Install]
WantedBy=default.target
